# examples/simpleB.service
[Unit]
Description=Foo
Before=simpleA

[Service]
ExecStart=/bin/echo B

[Install]
WantedBy=multi-user.target