# examples/units/simpleB.service
[Unit]
Description=Foo

[Service]
ExecStart=/packages/example_service 1

[Install]
WantedBy=multi-user.target