# examples/units/simpleA.service
[Unit]
Description=Foo
Before=simpleB

[Service]
ExecStart=/packages/example_service 10

[Install]
WantedBy=multi-user.target